# services/tor.1/run
#!/bin/sh
exec chpst -u tor:nogroup tor -f /etc/tor/torrc.1
